# systemd/bluetooth-audio.service
[Unit]
Description=Bluetooth A2DP Sink (bluealsad – NAP)
Documentation=https://github.com/arkq/bluez-alsa

PartOf=audio-bluetooth.target
After=sound.target bluetooth.service dbus.service audio-bluetooth.target
Requires=bluetooth.service dbus.service
Wants=bluetooth.service

[Service]
Type=dbus
BusName=org.bluealsa

# bluealsad exposes Bluetooth audio devices as ALSA PCM devices.
# Built from source (github.com/arkq/bluez-alsa) – not available in
# Raspberry Pi OS Bookworm default repos.
#
# -S               : log to syslog instead of stderr
# -p a2dp-sink     : accept A2DP audio pushed from a phone/tablet
# -p hfp-hf        : hands-free profile; remove if not needed
# Codecs advertised: SBC (mandatory) + FastStream (built-in, royalty-free)
ExecStart=/usr/bin/bluealsad -S \
    -p a2dp-sink \
    -p hfp-hf

# bluealsa-aplay captures the A2DP stream and feeds it to the default ALSA device.
# 00:00:00:00:00:00 = accept audio from any paired device.
ExecStartPost=/usr/bin/bluealsa-aplay --profile-a2dp 00:00:00:00:00:00

Restart=on-failure
RestartSec=3

# D-Bus access required; must not run in a private namespace that hides the bus.
PrivateTmp=no
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

[Install]
WantedBy=audio-bluetooth.target
