# systemd/shairport-sync.service
[Unit]
Description=AirPlay Audio Receiver (shairport-sync – NAP)
Documentation=https://github.com/mikebrady/shairport-sync

PartOf=audio-airplay.target
After=sound.target network-online.target avahi-daemon.service audio-airplay.target
Wants=avahi-daemon.service

[Service]
Type=simple
User=shairport-sync
Group=shairport-sync
SupplementaryGroups=audio

ExecStart=/usr/local/bin/shairport-sync \
    --output=alsa \
    --on-start=/bin/true \
    --on-stop=/bin/true

Restart=on-failure
RestartSec=3

# Prevent shairport-sync from touching anything outside its scope
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/run/shairport-sync

[Install]
WantedBy=audio-airplay.target
